# supervisord.conf
[supervisord]
nodaemon=true

[program:vaa3d-jobs]
command=/bin/bash -c "$TASK_RUNNER_PATH jobs"

[program:vaa3d-job-items]
command=/bin/bash -c "$TASK_RUNNER_PATH job_items"